# docker/cron/crontab
0 4 */8 * * /scripts/restart_container.sh
